# debian/postinst
#!/bin/sh

set -e

case "$1" in 
    configure)
    systemctl enable deepin-cd-helper.service
    ;;
esac


#DEBHELPER#
