# debian/prerm
#!/bin/sh

set -e

systemctl disable deepin-cd-helper.service

